# lxd/propane.service
[Unit]
Description=propane-service
After=network.target

[Service]
ExecStart=/app/server
RestartSec=50s

[Install]
WantedBy=multi-user.target
